# ansible/templates/s3uploader-client.service.j2
# Managed by Ansible — do not edit manually. Change in repo and re-deploy.
[Unit]
Description=S3 Uploader Client
After=network.target

[Service]
Type=simple
ExecStart={{ client_dir }}/s3up --config {{ client_dir }}/client.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
